# ai-sdlc-frontend.service
[Unit]
Description=AI SDLC Frontend
After=network.target

[Service]
Type=simple
User=ali-riaz
WorkingDirectory=/home/ali-riaz/Desktop/ai_sdlc/frontend
Environment=NODE_ENV=production
ExecStart=/home/ali-riaz/.nvm/versions/node/v20.20.2/bin/node /home/ali-riaz/Desktop/ai_sdlc/frontend/node_modules/.bin/next start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
